# recipes-bsp/autoinstall/files/autoinstall.service
[Unit]
Description=Install swu image

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/local/bin/install.sh

